# com.rize-to-calendar.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rize-to-calendar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/rize-to-calendar/.venv/bin/python</string>
        <string>/path/to/rize-to-calendar/rize_to_calendar.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/rize-to-calendar</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/rize-to-calendar/logs/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/rize-to-calendar/logs/sync.log</string>
</dict>
</plist>
